# roles/borgbackup_client/files/borgbackup-panic.service
[Unit]
Description=Notifies on failed %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borgbackup-panic %i